# scripts/crontab
0 * * * * root ~/update_cache.sh >> /var/log/cron.log 2>&1
